# nodeflow-edge.service
[Unit]
Description=Nodeflow Edge IoT Gateway
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=root
WorkingDirectory=/opt/nodeflow-edge
ExecStart=/opt/nodeflow-edge/venv/bin/python -m nodeflow_edge.main --config /opt/nodeflow-edge/config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nodeflow-edge

# Environment
Environment=PYTHONUNBUFFERED=1

# Watchdog
WatchdogSec=120

[Install]
WantedBy=multi-user.target
